# deploy/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:redis]
command=redis-server /etc/redis/redis.conf
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/redis-stdout.log
stderr_logfile=/var/log/supervisor/redis-stderr.log
priority=10

[program:trading-app]
command=python -m uvicorn src.main:app --host 0.0.0.0 --port 8000 --workers 1
directory=/app
autostart=true
autorestart=true
startsecs=10
startretries=3
stdout_logfile=/var/log/supervisor/app-stdout.log
stderr_logfile=/var/log/supervisor/app-stderr.log
environment=PYTHONPATH="/app",PYTHONUNBUFFERED="1"
priority=20

[program:alembic-migrate]
command=python -m alembic upgrade head
directory=/app
autostart=true
autorestart=false
startsecs=0
exitcodes=0
priority=15
stdout_logfile=/var/log/supervisor/migrate-stdout.log
stderr_logfile=/var/log/supervisor/migrate-stderr.log
